# infrastructure/ansible/roles/fleet-agent/templates/orbit.service.j2
# Fleet Agent (Orbit) Systemd Service
# Generated by Ansible for {{ inventory_hostname }}

[Unit]
Description=Fleet osquery agent (Orbit)
Documentation=https://fleetdm.com/docs
After=network.target
Wants=network.target

[Service]
Type=simple
User=root
Group=root

# Service execution
ExecStart=/opt/orbit/bin/orbit/orbit --config {{ fleet_agent_config_dir }}/orbit.yml
ExecReload=/bin/kill -HUP $MAINPID

# Restart configuration
Restart=always
RestartSec=10
StartLimitInterval=60
StartLimitBurst=3

# Resource limits
MemoryLimit=512M
CPUQuota=50%

# Security settings
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths={{ fleet_agent_config_dir }} /var/log/orbit /tmp

# Environment
Environment=ORBIT_LOG_LEVEL={{ fleet_agent_log_level | default('info') }}
Environment=ORBIT_CONFIG_PATH={{ fleet_agent_config_dir }}/orbit.yml

# Working directory
WorkingDirectory={{ fleet_agent_config_dir }}

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=orbit

[Install]
WantedBy=multi-user.target